# aa_deploy/bigness-league.service
[Unit]
Description=Bigness League Discord Bot
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=bigness
Group=bigness
WorkingDirectory=/home/bigness/DiscordBot/BignessLeague-DiscordBot

EnvironmentFile=-/home/bigness/DiscordBot/BignessLeague-DiscordBot/.env.production
EnvironmentFile=-/home/bigness/DiscordBot/BignessLeague-DiscordBot/.env

Environment=PYTHONPATH=/home/bigness/DiscordBot/BignessLeague-DiscordBot/src
Environment=PYTHONUNBUFFERED=1

ExecStart=/home/bigness/DiscordBot/BignessLeague-DiscordBot/.venv/bin/python -m bigness_league_bot.main
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target